# deploy/niravan.service
[Unit]
Description=NIRAVAN Cybersecurity Intelligence Platform Backend Service
After=network.target

[Service]
User=www-data
WorkingDirectory=/var/www/niravan
Environment="PATH=/var/www/niravan/venv/bin"
EnvironmentFile=/var/www/niravan/backend/.env
ExecStart=/var/www/niravan/venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000 --workers 4
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=niravan-api

[Install]
WantedBy=multi-user.target
